# rootfs/etc/s6-overlay/s6-rc.d/init-validate/run
#!/command/with-contenv bash
# shellcheck shell=bash
# shellcheck source=/usr/local/lib/s6-functions.sh

set -euo pipefail
IFS=$'\n\t'

. /usr/local/lib/s6-functions.sh
_load_s6_vars

#---------------------------------------------------------------------------------------------------------------------------
# HELPERS
#---------------------------------------------------------------------------------------------------------------------------
function _require_var() {
    local name="$1"
    local value="${!name:-}"

    if [[ -z "${value}" ]]; then
        _error "${CURRENT_DIR} | ${name} is not set"
        return 1
    fi
    return 0
}

function _warn_if_default() {
    local name="$1"
    local bad="$2"
    local value="${!name:-}"

    if [[ -n "${value}" && "${value}" == "${bad}" ]]; then
        _warn "${CURRENT_DIR} | Default ${name} detected, this is a security risk, please set a secure value"
    fi
}

#---------------------------------------------------------------------------------------------------------------------------
# VALIDATION: PLATFORM / MODE
#---------------------------------------------------------------------------------------------------------------------------
function _init_validate_y2038_issue() {
    local year arch
    year="$(date +%Y)"
    arch="$(uname -m)"

    if [[ "${arch}" == "armv7l" && ( "${year}" == "1969" || "${year}" == "1970" ) ]]; then
        _error "${CURRENT_DIR} | Y2038 issue detected: Your Docker host has an outdated libseccomp version. Please update libseccomp on your host system"
        return 1
    fi
    return 0
}

function _init_validate_operation_mode() {
    local mode="${CONTAINER_RUN_MODE:-standalone}"

    case "${mode}" in
        standalone|server_remote_datalogger|remote_datalogger)
            _info "${CURRENT_DIR} | Operation mode: ${mode}"
            return 0
            ;;
        *)
            _error "${CURRENT_DIR} | Invalid CONTAINER_RUN_MODE: ${mode}"
            _error "${CURRENT_DIR} | Valid options: standalone, server_remote_datalogger, remote_datalogger"
            return 1
            ;;
    esac
}

#---------------------------------------------------------------------------------------------------------------------------
# VALIDATION: REMOTE DATALOGGER SUBCONFIG
#---------------------------------------------------------------------------------------------------------------------------
function _init_datalogger_mode_validate_serial_config() {
    local errors=0

    _require_var "DSMRREADER_REMOTE_DATALOGGER_SERIAL_DEVICE" || ((errors++))
    _require_var "DSMRREADER_REMOTE_DATALOGGER_SERIAL_BAUDRATE" || ((errors++))
    _require_var "DSMRREADER_REMOTE_DATALOGGER_SERIAL_BYTESIZE" || ((errors++))

    [[ "${errors}" -eq 0 ]]
}

function _init_datalogger_mode_validate_network_config() {
    local errors=0

    _require_var "DSMRREADER_REMOTE_DATALOGGER_NETWORK_HOST" || ((errors++))
    _require_var "DSMRREADER_REMOTE_DATALOGGER_NETWORK_PORT" || ((errors++))

    [[ "${errors}" -eq 0 ]]
}

function _init_datalogger_mode_validate_remote_datalogger_config() {
    local errors=0

    _require_var "DSMRREADER_REMOTE_DATALOGGER_API_HOSTS" || ((errors++))
    _require_var "DSMRREADER_REMOTE_DATALOGGER_API_KEYS" || ((errors++))

    [[ "${errors}" -eq 0 ]]
}

#---------------------------------------------------------------------------------------------------------------------------
# VALIDATION: ILLEGAL VAR COMBINATIONS
#---------------------------------------------------------------------------------------------------------------------------
function _init_validate_remote_vars_vs_mode() {
    local mode="${CONTAINER_RUN_MODE:-standalone}"
    local found=()

    while IFS='=' read -r k v; do
        [[ -z "${v}" ]] && continue
        [[ "${k}" == DSMRREADER_REMOTE_DATALOGGER_* ]] && found+=("${k}")
    done < <(env)

    if [[ "${mode}" != "remote_datalogger" && "${#found[@]}" -gt 0 ]]; then
        _error "${CURRENT_DIR} | Remote datalogger variables detected, but CONTAINER_RUN_MODE is '${mode}'"
        _error "${CURRENT_DIR} | These variables are only valid when CONTAINER_RUN_MODE=remote_datalogger"
        _error "${CURRENT_DIR} | Found: ${found[*]}"
        return 1
    fi

    if [[ "${mode}" == "remote_datalogger" ]]; then
        if [[ -n "${DJANGO_DATABASE_HOST:-}${DJANGO_DATABASE_NAME:-}${DJANGO_DATABASE_USER:-}${DJANGO_DATABASE_PASSWORD:-}${DJANGO_DATABASE_PORT:-}" ]]; then
            _warn "${CURRENT_DIR} | Database variables are set but not used in remote_datalogger mode"
        fi
        if [[ -n "${DSMRREADER_ADMIN_USER:-}${DSMRREADER_ADMIN_PASSWORD:-}" ]]; then
            _warn "${CURRENT_DIR} | Admin variables are set but not used in remote_datalogger mode"
        fi
    fi

    return 0
}

#---------------------------------------------------------------------------------------------------------------------------
# VALIDATION: REQUIRED VARS (MODE BASED)
#---------------------------------------------------------------------------------------------------------------------------
function _init_validate_required_vars() {
    local errors=0
    local mode="${CONTAINER_RUN_MODE:-standalone}"

    case "${mode}" in
        standalone|server_remote_datalogger)
            _require_var "DSMRREADER_ADMIN_USER" || ((errors++))
            _require_var "DSMRREADER_ADMIN_PASSWORD" || ((errors++))
            _warn_if_default "DSMRREADER_ADMIN_PASSWORD" "admin"

            _require_var "DJANGO_SECRET_KEY" || ((errors++))
            _warn_if_default "DJANGO_SECRET_KEY" "dsmrreader"

            _require_var "DJANGO_DATABASE_NAME" || ((errors++))
            _require_var "DJANGO_DATABASE_HOST" || ((errors++))
            _require_var "DJANGO_DATABASE_USER" || ((errors++))
            _require_var "DJANGO_DATABASE_PASSWORD" || ((errors++))
            _warn_if_default "DJANGO_DATABASE_PASSWORD" "dsmrreader"
            _require_var "DJANGO_DATABASE_PORT" || ((errors++))
            ;;
        remote_datalogger)
            if ! _init_datalogger_mode_validate_remote_datalogger_config; then
                ((errors++))
            fi

            local input_method="${DSMRREADER_REMOTE_DATALOGGER_INPUT_METHOD:-serial}"

            case "${input_method}" in
                serial)
                    _info "${CURRENT_DIR} | Remote datalogger input method: serial"
                    _init_datalogger_mode_validate_serial_config || ((errors++))
                    ;;
                ipv4)
                    _info "${CURRENT_DIR} | Remote datalogger input method: ipv4"
                    _init_datalogger_mode_validate_network_config || ((errors++))
                    ;;
                *)
                    _error "${CURRENT_DIR} | Invalid DSMRREADER_REMOTE_DATALOGGER_INPUT_METHOD: ${input_method}"
                    ((errors++))
                    ;;
            esac
            ;;
        *)
            _error "${CURRENT_DIR} | Invalid CONTAINER_RUN_MODE: ${mode}"
            return 1
            ;;
    esac

    if [[ "${errors}" -gt 0 ]]; then
        _error "${CURRENT_DIR} | Missing or invalid configuration for mode: ${mode} (${errors} problem(s))"
        return 1
    fi

    return 0
}

#---------------------------------------------------------------------------------------------------------------------------
# MAIN EXECUTION
#---------------------------------------------------------------------------------------------------------------------------
function main() {
    [[ "${CONTAINER_ENABLE_DEBUG:-false}" == "true" ]] && set -x

    CURRENT_DIR="$(basename "$(dirname "$0")")"
    _info "${CURRENT_DIR} | Starting"

    _init_validate_y2038_issue || exit 1
    _init_validate_operation_mode || exit 1
    _init_validate_remote_vars_vs_mode || exit 1
    _init_validate_required_vars || exit 1

    _success "${CURRENT_DIR} | Validation complete"
}

main "$@"
