# rootfs/etc/s6-overlay/s6-rc.d/svc-dsmr-datalogger-log/run
#!/bin/execlineb -P
# /etc/s6-overlay/s6-rc.d/dsmr_datalogger-log/run
# shellcheck shell=bash

s6-log -bp -- /var/log/dsmr_datalogger
