# rootfs/etc/s6-overlay/s6-rc.d/svc-dsmr-remote-datalogger/run
#!/command/with-contenv bash
# shellcheck shell=bash
# shellcheck source=/usr/local/lib/s6-functions.sh

set -euo pipefail
IFS=$'\n\t'

. /usr/local/lib/s6-functions.sh
_load_s6_vars

svc_dsmr_remote_datalogger_start() {
    if [[ "${CONTAINER_RUN_MODE}" = "remote_datalogger" ]]; then
        CURRENT_DIR="$(basename "$(pwd)")"
        _info "${CURRENT_DIR} | Starting"

        cmd=(/opt/venv/bin/python -u /app/dsmr_datalogger/scripts/dsmr_datalogger_api_client.py)

        if command -v ionice >/dev/null 2>&1; then
            exec /command/s6-setuidgid app ionice -c 2 -n 7 nice -n 5 "${cmd[@]}"
        else
            exec /command/s6-setuidgid app nice -n 5 "${cmd[@]}"
        fi
    fi

    CURRENT_DIR="$(basename "$(pwd)")"
    _info "${CURRENT_DIR} | Disabled for mode ${CONTAINER_RUN_MODE}"
    exec sleep infinity
}

main() {
    [[ "${CONTAINER_ENABLE_DEBUG:-false}" == "true" ]] && set -x
    svc_dsmr_remote_datalogger_start
}

main "$@"
