# src/root/etc/services.d/dsmr_datalogger/run
#!/usr/bin/with-contenv bash

if [[ "${DATALOGGER_MODE}" = standalone ]]; then
    echo "Starting DSMR Datalogger..."
    /usr/local/bin/python3 -u /dsmr/dsmr_datalogger_api_client.py
fi
