# rootfs/etc/s6-overlay/s6-rc.d/svc-dsmr-backend/run
#!/command/with-contenv bash
# shellcheck shell=bash

function _info()  { printf "\r[ \033[00;34mINFO\033[0m ] %s\n" "$@"; }

if [[ "${DSMRREADER_OPERATION_MODE}" = standalone || "${DSMRREADER_OPERATION_MODE}" = api_server ]]; then
    _info "Starting DSMR Reader - backend..."
    cd /app || exit
    exec s6-setuidgid app /opt/venv/bin/python -u /app/manage.py dsmr_backend
else
    _info "svc-dsmr-backend disabled for mode ${DSMRREADER_OPERATION_MODE}"
    sleep infinity
fi
