# rootfs/etc/s6-overlay/s6-rc.d/svc-nginx/run
# shellcheck disable=SC1008
#!/command/with-contenv bash

exec 2>&1

if [[ "${DSMRREADER_OPERATION_MODE}" = standalone || "${DSMRREADER_OPERATION_MODE}" = api_server ]]; then
    echo "Starting DSMR Reader - nginx..."
    cd /app || exit
    exec /usr/sbin/nginx -g 'daemon off;'
else
    sleep infinity
fi
