# rootfs/etc/s6-overlay/s6-rc.d/svc-nginx/run
#!/command/with-contenv bash
# shellcheck shell=bash

if [[ "${DSMRREADER_OPERATION_MODE}" = standalone || "${DSMRREADER_OPERATION_MODE}" = api_server ]]; then
    for i in {1..10}; do
        if [ -e "/tmp/gunicorn--dsmr_webinterface.socket" ]; then
            echo "Starting DSMR Reader - nginx..."
            sleep 2
            cd /app || exit
            /usr/sbin/nginx -g 'daemon off;'
        else
            sleep 1
            printf "\\rTesting webservice socket connectivity: %s second(s) of 10 seconds..." "$i"
            if [[ $i == 10 ]]; then
                echo "Webinterface socket couldn't be verified! Please verify the webinterface log files. Exiting..."
                exit 1
            fi
        fi
    done
else
    sleep infinity
fi
