# rootfs/etc/s6-overlay/s6-rc.d/svc-nginx/run
#!/command/with-contenv bash
# shellcheck shell=bash
# shellcheck source=/usr/local/lib/s6-functions.sh

. /usr/local/lib/s6-functions.sh
_load_s6_vars

function svc_nginx_start() {
    if [[ "${DSMRREADER_CONTAINER_RUN_MODE}" = standalone || "${DSMRREADER_CONTAINER_RUN_MODE}" = api_server ]]; then
        for i in {1..10}; do
            if [[ -e "/run/gunicorn/dsmr_webinterface.socket" ]]; then
                _success "${CURRENT_DIR} webinterface socket ready"
                sleep 2

                exec /usr/sbin/nginx -g 'daemon off;'
            fi

            _info "${CURRENT_DIR} is waiting for webinterface socket... (${i}/10)"
            sleep 1
        done

        _error "${CURRENT_DIR} webinterface socket not found after 10 attempts"
        exit 1
    else
        _info "${CURRENT_DIR} disabled for mode ${DSMRREADER_CONTAINER_RUN_MODE}"
        exit 0
    fi
}

function main() {
    [[ "${DSMRREADER_CONTAINER_DEBUG}" == "true" ]] && set -x

    CURRENT_DIR="$(basename "$(pwd)")"
    svc_nginx_start
}

main "$@"