# rootfs/etc/s6-overlay/s6-rc.d/svc-nginx/run
#!/bin/execlineb -P
# /etc/s6-overlay/s6-rc.d/nginx/run
# shellcheck shell=bash

with-contenv

if [[ "${DSMRREADER_OPERATION_MODE}" = standalone || "${DSMRREADER_OPERATION_MODE}" = api_server ]]; then
    echo "Starting DSMR Reader - nginx..."
    cd /app || exit
    s6-setuidgid app
    /usr/sbin/nginx -g 'daemon off;'
else
    sleep infinity
fi
